# layout/Library/LaunchDaemons/XcodeAnyDebug-startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>ExecuteAllowed</key>
	<true/>
	<key>Label</key>
	<string>XcodeAnyDebug.startup</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>@JBROOT@/bin/sh</string>
	<key>ProgramArguments</key>
	<array>
		<string>@JBROOT@/bin/sh</string>
		<string>@JBROOT@/usr/libexec/XcodeAnyDebug-startup</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>